# adb.conf
[program:adb]
priority=10
user=root
command=bash /root/adb.sh
